# s6/services.d/anti-timeout/run
#!/bin/sh

MSG="<!--anti-timeout--!>"
while true; do
	echo "$MSG"
	sleep 60
done
